# ansible/roles/common/templates/kolla.target.j2
[Unit]
Description=Kolla target allowing to start/stop all kolla*@.service instances at once

[Install]
WantedBy=multi-user.target
