# tools/virt-guest-shutdown.target
[Unit]
Description=Coordination target for machined controlled guests

[Install]
WantedBy=multi-user.target
